# config/systemd/puma.socket.erb
[Unit]
Description = <%= fetch(:application) %>'s Puma HTTP Server Accept Sockets

[Socket]
ListenStream = <%= current_path %>/tmp/sockets/puma.sock
NoDelay = true
ReusePort = true
Backlog = 1024

[Install]
WantedBy=sockets.target
