# config/systemd/sneakers.service.erb
[Unit]
Description = <%= fetch(:application) %> background jobs service

[Service]
Environment = RAILS_ENV=<%= fetch(:rails_env) %>
Environment = PWD=<%= current_path %>
WorkingDirectory = <%= current_path %>
PIDFile = <%= shared_path %>/tmp/pids/sneakers.pid

ExecStart = <%= SSHKit.config.command_map[:rvm] %> default do bundle exec rake sneakers:run
ExecReload = /bin/kill -SIGUSR1 $MAINPID
ExecStop = /bin/kill -SIGTERM $MAINPID
User = <%= ENV["STAGING_USER"] %>
Group = <%= ENV["STAGING_USER"] %>
Restart = always

[Install]
WantedBy = multi-user.target
